# deploy/cloudflared.service
[Unit]
Description=Cloudflare Tunnel — seginternet.rochainf.com.br
Documentation=https://developers.cloudflare.com/cloudflare-one/connections/connect-apps/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/cloudflared tunnel --config /etc/cloudflared/config.yml run
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloudflared

# Reiniciar em caso de falha com backoff progressivo
StartLimitIntervalSec=60
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
